# update-ds4-leds.service
[Unit]
Description=Updates DualShock 4 controller LED to primary color of Noctalia Shell
[Service]
Type=oneshot
ExecStart=/usr/bin/update-ds4-leds